# examples/ModBus/deb/modbusex.service
[Unit]
Description=ModBus Example

[Service]
Type=notify
WorkingDirectory=/opt/ModBus
ExecStart=/opt/ModBus/ModBus

[Install]
WantedBy=multi-user.target