# examples/TemperatureController/tempcontroller.service
[Unit]
Description=Temperature Controller

[Service]
Type=notify
ExecStart=/opt/TemperatureController/TemperatureController

[Install]
WantedBy=multi-user.target